# deploy/systemd/hohai-checkin.service
[Unit]
Description=Hohai Daily Check-in (Python)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/HohaiCheckin
EnvironmentFile=-/etc/hohai-checkin.env
ExecStart=/usr/bin/python3 /opt/HohaiCheckin/checkin.py
User=actions
Group=actions
Nice=10
StandardOutput=append:/var/log/hohai-checkin.log
StandardError=append:/var/log/hohai-checkin.log
