# deploy/systemd/hohai-checkin-docker.service
[Unit]
Description=Hohai Check-in (Docker Compose)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/HohaiCheckin
RemainAfterExit=yes
ExecStart=/usr/bin/docker compose run --rm hohai-checkin
ExecStop=/bin/true

[Install]
WantedBy=multi-user.target
